# crontab
* * * * * root /usr/local/bin/node /app/worker.js >> /tmp/cron.log 2>&1
